# bql-enterprise/deploy/systemd/bqs-enterprise.service
[Unit]
Description=BQS Enterprise Adapter Service (template)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=bqs
Group=bqs

# Template only — wire to your enterprise binary when available.
# ExecStart=/opt/bqs/bin/bqs-enterprise --config /etc/bqs/config.yaml

Restart=on-failure
RestartSec=2
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true

[Install]
WantedBy=multi-user.target
